# xiboplayer/xiboplayer-chromium.service
[Unit]
Description=Xibo Digital Signage Player (Chromium Kiosk)
Documentation=https://github.com/xibo-players/xiboplayer-chromium
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/xiboplayer-chromium
Restart=always
RestartSec=5
TimeoutStopSec=15
TimeoutStopFailureMode=terminate
KillMode=mixed
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
LimitCORE=0
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xiboplayer-chromium

[Install]
WantedBy=graphical-session.target
